# sonarr.ini
[program:sonarr]
user = sonarr
group = sonarr
command = /usr/bin/sonarr -nobrowser -data=/var/lib/sonarr
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
